# crontab
0 4 * * * certbot renew --post-hook /update-haproxy-certs.sh

